# deploy/podman/suse-ai-ingest.timer
# =============================================================================
# suse-ai-ingest.timer — nightly timer for document re-indexing
# Install: cp suse-ai-ingest.timer ~/.config/systemd/user/
# Enable:  systemctl --user enable --now suse-ai-ingest.timer
# =============================================================================

[Unit]
Description=SUSE AI Nightly Document Ingestion Timer
Documentation=https://github.com/suse/suse-ai
After=network-online.target

[Timer]
OnCalendar=*-*-* 03:00:00
Persistent=true
RandomizedDelaySec=300
AccuracySec=1m

[Install]
WantedBy=timers.target
