# deploy/podman/suse-ai.container
# =============================================================================
# suse-ai.container — Quadlet container definition for Podman 5.x
# Install: cp suse-ai.container ~/.config/containers/systemd/
#          cp suse-ai-pod.pod   ~/.config/containers/systemd/
# Reload:  systemctl --user daemon-reload
# Start:   systemctl --user enable --now suse-ai-pod suse-ai
# Logs:    journalctl --user -u suse-ai -f
#
# 2026 FIXES:
#   - Cgroupsmode=no-conmon  REQUIRED when joining a pod (prevents cgroup conflict)
#   - UserNS=keep-id         Required for rootless volume ownership
#   - HealthCmd port fixed to 8090 (was wrong 8080 in original)
#   - Do NOT create suse-ai.service alongside this file — Quadlet auto-generates
#     it. Having both causes a unit name collision.
# =============================================================================

[Container]
Image=localhost/suse-ai:latest
Pod=suse-ai-pod
Volume=/var/lib/suse-ai:/data:rw,z
Environment=PODMAN_USERNS=keep-id
Environment=LLM_BASE_URL=http://llm:11434
Environment=EMBED_BASE_URL=http://embed:11435
Environment=APP_PORT=8090
Environment=LOG_LEVEL=info
UserNS=keep-id
Cgroupsmode=no-conmon
AutoUpdate=registry
HealthCmd=curl -sf http://localhost:8090/health || exit 1
HealthInterval=30s
HealthTimeout=10s
HealthStartPeriod=60s
HealthRetries=3
Label=io.podman.version=5.4.2
Label=io.opensuse.leap=16
Label=app=suse-ai
NoNewPrivileges=true
SecurityLabelDisable=false
LogDriver=journald

[Service]
Restart=on-failure
RestartSec=10
MemoryMax=3G
TimeoutStopSec=30

[Install]
WantedBy=default.target
