# deploy/podman/suse-ai.socket
# =============================================================================
# suse-ai.socket — systemd socket unit for SUSE AI application
# Install: cp suse-ai.socket ~/.config/systemd/user/
# Enable:  systemctl --user enable --now suse-ai.socket
# =============================================================================

[Unit]
Description=SUSE AI Application Socket (port 8090)
Documentation=https://github.com/suse/suse-ai

[Socket]
ListenStream=8090
Accept=no
NoDelay=true
KeepAlive=true
KeepAliveTimeSec=30
KeepAliveIntervalSec=5
KeepAliveProbes=3
Backlog=128

[Install]
WantedBy=default.target
